# templates/program.conf.erb
[program:<%= @app %>]
command=/usr/bin/php index.php worker
numprocs=<%= @number_worker %>
process_name=%(program_name)s_%(process_num)02d
directory=<%= @working_dir %>
stdout_logfile=/var/log/gearman-job-server/<%= @app %>.log
environment=GEARMAN_USER=gearman,CI_ENV=production
autostart=true
autorestart=true
user=www-data
stopsignal=KILL
